# deploy/systemd/wormhole-nightly-flush.timer
[Unit]
Description=Run Wormhole nightly queue flush at midnight Pacific

[Timer]
OnCalendar=*-*-* 00:00:00
Timezone=America/Los_Angeles
Persistent=true
Unit=wormhole-nightly-flush.service

[Install]
WantedBy=timers.target